# deploy/com.newsclassified.marketstream.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.newsclassified.marketstream</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>/Users/admin/Library/Application Support/news_classified/scripts/run_local.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/admin/Library/Application Support/news_classified</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
      <key>MARKET_STREAM_ACCESS_PASSWORD</key>
      <string>123456</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/admin/Library/Application Support/news_classified/logs/market_stream.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/admin/Library/Application Support/news_classified/logs/market_stream.err.log</string>
  </dict>
</plist>
